# btc_bot/daily.cron
* * * * * python3 /run.py push >> /tmp/cron.log 2>&1
