# btc_bot/daily.cron
* * * * * echo "✅ Cron ran at $(date)" >> /tmp/cron_test.log 2>&1
